# bash/crontabs/backup
* 20 * * * /Users/garethjones/.bash/bin/backup_to_ephalent.sh
